# RG351P-M/Ubuntu OS Partition/var/lib/dpkg/info/retroroller.list
/.
/etc
/etc/systemd
/etc/systemd/system
/etc/systemd/system/multi-user.target.wants
/home
/home/odroid
/home/odroid/.config
/home/odroid/.config/retroarch
/home/odroid/.config/retroarch/retroarch.cfg
/home/odroid/.config/retroarch32
/home/odroid/.config/retroarch32/retroarch.cfg
/lib
/lib/systemd
/lib/systemd/system-sleep
/lib/systemd/system-sleep/sleep
/opt
/opt/retroarch
/opt/retroarch/bin
/opt/retroarch/bin/retroarch
/opt/retroarch/bin/retroarch32
/opt/retroroller
/opt/retroroller/cmds
/opt/retroroller/cmds/retroarch.sh
/opt/retroroller/cmds/retroarch32.sh
/opt/retroroller/scripts
/opt/retroroller/scripts/retroroller_boot.sh
/opt/retroroller/services
/opt/retroroller/services/oga_events.service
/opt/retroroller/services/retroroller_boot.service
/usr
/usr/local
/usr/local/bin
/usr/local/bin/oga_events.py
/usr/local/bin/retroarch
/usr/local/bin/retroarch32
/usr/local/lib
/usr/local/lib/aarch64-linux-gnu
/usr/local/lib/aarch64-linux-gnu/librga.so
/usr/local/lib/arm-linux-gnueabihf
/usr/local/lib/arm-linux-gnueabihf/librga.so
/etc/systemd/system/multi-user.target.wants/alsa-state.service
/etc/systemd/system/multi-user.target.wants/oga_events.service
/etc/systemd/system/multi-user.target.wants/retroroller_boot.service
/etc/systemd/system/oga_events.service
/etc/systemd/system/retroroller_boot.service
