# RG351P-M/Ubuntu OS Partition/var/lib/dpkg/info/utils-go2_arm64.preinst
#!/bin/bash
deb-systemd-invoke stop audiopath.service >/dev/null || true
systemctl disable audiopath.service >/dev/null || true
