# RG351P-M/Ubuntu OS Partition/var/lib/dpkg/info/utils-go2_arm64.postinst
#!/bin/bash
systemctl enable audiopath.service >/dev/null || true
deb-systemd-invoke start audiopath.service >/dev/null || true
