# RG351P-M/Ubuntu OS Partition/var/lib/dpkg/info/utils-go2_arm64.list
/.
/etc
/etc/systemd
/etc/systemd/system
/etc/systemd/system/audiopath.service
/opt
/opt/system
/opt/system/Change Password.sh
/opt/system/Network Info.sh
/opt/system/Wifi.sh
/usr
/usr/bin
/usr/bin/audiopath
/usr/bin/bioscheck
/usr/bin/msgbox
/usr/bin/nmui
/usr/bin/osk
