# RG351P-M/Ubuntu OS Partition/etc/systemd/system/pico8hotkey.service
[Unit]
Description=Pico 8 Exit Hotkey service

StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1

ExecStart=/usr/local/bin/pico8keydemon.py

[Install]
WantedBy=multi-user.target

